# podman-refactor/plugin-bigwig.container
[Unit]
Description=plugin-bigwig

[Container]
# Use the centos image
Image=docker.io/veupathdb/vdi-plugin-bigwig:latest
ContainerName=${COMPOSE_PROJECT_NAME}-%N

# Use volume and network defined below
Network=internal.network
NetworkAlias=bigwig
#-monitoring-ext: {}

EnvironmentFile=env/env_common_defaults
EnvironmentFile=env/env_plugin_bigwig_defaults
Environment=DATASET_INSTALL_ROOT=${DATASET_DIRECTORY_TARGET_PATH}

Secret=DB_CONNECTION_PASS_AMOEBA,type=env

#Volume=${DATASET_DIRECTORY_SOURCE_PATH}:${DATASET_DIRECTORY_TARGET_PATH}:Z
Volume=datasets.volume:${DATASET_DIRECTORY_TARGET_PATH}:Z

#Replaces watchtower with podman-auto-update
AutoUpdate=registry

#    labels:
#    - "prometheus.scrape_enabled=true"
[Service]
EnvironmentFile=/etc/containers/systemd/vdi/env/env_common_defaults
